# install/battery-charge-limit.service
[Unit]
Description=Set battery charge limit to 60%%
ConditionPathExistsGlob=/sys/class/power_supply/BAT*/charge_control_end_threshold

[Service]
Type=oneshot
ExecStart=/bin/sh -c 'for battery in /sys/class/power_supply/BAT*; do [ -e "$battery/charge_control_end_threshold" ] && echo 60 > "$battery/charge_control_end_threshold"; done'
